# aws/boot/prefer-boot.service
[Unit]
Description=PreFer llama.cpp inference (stage NVMe + run container)
# Needs the Docker daemon, outbound network (S3 / Hugging Face / registry), and
# the DLAMI NVMe mount (dlami-nvme formats/mounts the instance store on boot).
After=docker.service network-online.target dlami-nvme.service
Wants=network-online.target dlami-nvme.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/opt/prefer/prefer-boot.env
ExecStart=/opt/prefer/10-prep-nvme.sh
ExecStart=/opt/prefer/20-run-container.sh
ExecStop=/usr/bin/docker stop prefer
# First boot downloads tens of GB; never time the startup out.
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
